# instancesetup/crontab
*/5 * * * * /home/ubuntu/OnDemandMinecraft/autoshutdown.sh
